# deploy/vcare/supervisor/conf.d/nginx-error-backend.conf
[program:nginx-error-backend]
command=tail -f /var/log/nginx/backend_error.log
stdout_events_enabled=true
stderr_events_enabled=true
