# deploy/vcare/supervisor/conf.d/nginx-error.conf
[program:nginx-error]
command=tail -f /var/log/nginx/error.log
stdout_events_enabled=true
stderr_events_enabled=true
