# deploy/vcare/supervisor/conf.d/php-fpm-log.conf
[program:php-fpm-log]
command=tail -f /var/log/php7.0-fpm.log
stdout_events_enabled=true
stderr_events_enabled=true
