# deploy/vcare/supervisor/conf.d/toro-news-log.conf
[program:toro-news-log]
command=tail -f /var/www/shared/app/logs/toro.news.publish.log
stdout_events_enabled=true
stderr_events_enabled=true
