# deploy/vcare/supervisor/conf.d/toro-news-log-error.conf
[program:toro-news-log-error]
command=tail -f /var/www/shared/app/logs/toro.news.publish.error
stdout_events_enabled=true
stderr_events_enabled=true
